# target/linux/layerscape/base-files/etc/init.d/boot
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=10
STOP=98

uci_apply_defaults() {
	. /lib/functions/system.sh

	cd /etc/uci-defaults || return 0
	files="$(ls)"
	[ -z "$files" ] && return 0
	mkdir -p /tmp/.uci
	for file in $files; do
		( . "./$(basename $file)" ) && rm -f "$file"
	done
	uci commit
}

set_interface_name(){
	ip link set eth1 name iface0
	ip link set eth2 name iface1
	ip link set eth4 name iface2
	ip link set eth0 name iface3
	ip link set eth3 name iface4
}

check_update_etc(){
	# add by liji to support 检查是否需要恢复配置
	mkdir /recover && mount /dev/mmcblk0p3 /recover
	if [  -f /recover/sys_etc_recover ] ; then
		sysupgrade -r /recover/backup.tar.gz
		rm /recover/sys_etc_recover && sync	
		reboot		
	fi
	
	umount /recover && rm -rf /recover
}

set_led_status(){
	echo 0 > /sys/class/leds/cloud/delay_on
	echo 1000 > /sys/class/leds/cloud/delay_off
	
	echo 1000 > /sys/class/leds/sys/delay_on
	echo 1000 > /sys/class/leds/sys/delay_off
}

boot() {
	[ -f /proc/mounts ] || /sbin/mount_root
	[ -f /proc/jffs2_bbc ] && echo "S" > /proc/jffs2_bbc
	[ -f /proc/net/vlan/config ] && vconfig set_name_type DEV_PLUS_VID_NO_PAD

	mkdir -p /var/run
	mkdir -p /var/log
	mkdir -p /var/lock
	mkdir -p /var/state
	mkdir -p /var/tmp
	mkdir -p /tmp/.uci
	chmod 0700 /tmp/.uci
	touch /var/log/wtmp
	touch /var/log/lastlog
	touch /tmp/resolv.conf.auto
	ln -sf /tmp/resolv.conf.auto /tmp/resolv.conf
	grep -q debugfs /proc/filesystems && /bin/mount -o noatime -t debugfs debugfs /sys/kernel/debug
	[ "$FAILSAFE" = "true" ] && touch /tmp/.failsafe

	/sbin/kmodloader

	[ ! -f /etc/config/wireless ] && {
		# compat for brcm47xx and mvebu
		sleep 1
	}

	/bin/config_generate
	uci_apply_defaults
	
	# temporary hack until configd exists
	/sbin/reload_config
	# set wifi country by liji
	iw reg set CN
	
	#检查是否需要恢复配置
	check_update_etc
	#重置led状态
	set_led_status
	#set name interface
	set_interface_name
	#add check update command by liji
	/sbin/reset_gpio.sh &
	/sbin/process_check.sh &
	
}
