# itb/openwrt-layerscape-armv8_64b-device-ls1043ardb-rootfs/usr/lib/opkg/info/ipsec-tools.list
/usr/sbin/setkey
/usr/lib/libracoon.so.0
/etc/racoon/functions.sh
/usr/lib/libracoon.so.0.0.0
/usr/sbin/racoon
/etc/init.d/racoon
/usr/lib/libipsec.so.0
/etc/racoon/vpnctl
/etc/racoon/p1client-up
/etc/config/racoon
/usr/lib/libipsec.so.0.0.1
/etc/racoon/p1client-down
/usr/sbin/plainrsa-gen
/usr/sbin/racoonctl
