# itb/openwrt-layerscape-armv8_64b-device-ls1043ardb-rootfs/usr/lib/opkg/info/strongswan.list
/etc/init.d/ipsec
/usr/lib/ipsec/libstrongswan.so.0
/usr/lib/ipsec/libstrongswan.so.0.0.0
/etc/ipsec.conf
/etc/ipsec.secrets
/etc/strongswan.conf
/lib/upgrade/keep.d/strongswan
/etc/ipsec.user
